# dist/seaweedfs@.service
[Unit]
Description=SeaweedFS Service (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=all
User=root
Group=root
ExecStart=/opt/seaweedfs/seaweedfs-service.sh %i
Restart=always
RestartSec=5s
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
